# AppDir/usr/share/darktable/mcp/dt-mcp.service
[Unit]
Description=Darktable MCP AI Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.config/darktable/mcp/dt_mcp_server.py
WorkingDirectory=%h/.config/darktable/mcp
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
